# mkosi-initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup@.service.d/credential.conf
[Service]
ImportCredential=cryptsetup.*
# Compat with older systemd versions that don't support ImportCredential=.
LoadCredential=cryptsetup.passphrase
LoadCredential=cryptsetup.fido2-pin
LoadCredential=cryptsetup.tpm2-pin
LoadCredential=cryptsetup.luks2-pin
LoadCredential=cryptsetup.pkcs11-pin
